# sigma_auto_final/process_creation/132_detects_the_use_of_regexe_to_add_a_value_to_the_hkcusoftware.yml
id: 2835ec1d-97fd-4c24-a482-d5527e9927fb
status: test
description: Detects the use of reg.exe to add a value to the HKCU\Software\Microsoft\Windows\CurrentVersion\Run
  key.
author: AutoGenerated
date: '2025-07-10'
modified: 2022-11-29
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
  - Image|endswith: \reg.exe
  - OriginalFileName: reg.exe
  selection_cli:
    CommandLine|contains|all:
    - add
    - HKCU\Software\Microsoft\Windows\CurrentVersion\Run
  condition: all of selection_*
falsepositives:
- Unknown
level: medium
title: '[AUTO] Detects the use of reg.exe to add a value to the HKCU\Software\Microsoft\Windows\CurrentVersion\Run
  key.'
