# autogen_sigma/process_creation/014_user_generate_a_detection_rule_for_detect_scheduled_task_cre.yml
author: AutoGenerated
date: '2025-07-10'
modified: 2022-12-25
tags:
- attack.persistence
- attack.t1053.005
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
  - Image|endswith: \schtasks.exe
  - OriginalFileName: schtasks.exe
  selection_cli:
    CommandLine|contains:
    - ' /create '
    - ' /change '
  condition: all of selection_*
falsepositives:
- Legitimate use of scheduled tasks
level: medium
id: f501f196-15da-42a4-b1a7-dbd100aa1672
title: '[AUTO] User: Generate a detection rule for Detect scheduled task creation
  for persistence.'
status: experimental
description: 'User: Generate a detection rule for Detect scheduled task creation for
  persistence.'
