# autogen_sigma/process_creation/039_detects_suspicious_atexe_scheduled_tasks.yml
id: 4e479bc5-de7d-4937-ab18-2650d844f30d
status: test
description: Detects suspicious at.exe scheduled tasks
author: AutoGenerated
date: '2025-07-10'
modified: 2023-01-02
tags:
- attack.execution
- attack.t1053.005
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
  - Image|endswith: \at.exe
  - OriginalFileName: at.exe
  selection_cli:
    CommandLine|contains:
    - ' /every:M'
    - ' /next:'
    - ' /create'
  condition: all of selection_*
falsepositives:
- Unknown
level: medium
title: '[AUTO] Detects suspicious at.exe scheduled tasks'
