# autogen_sigma/process_creation/171_detects_the_use_of_regexe_to_modify_the_run_registry_key.yml
id: 4477f351-4fd2-4c1f-aa09-431375c7d10d
status: test
description: Detects the use of reg.exe to modify the "Run" registry key.
author: AutoGenerated
date: '2025-07-10'
modified: 2022-11-29
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
  - Image|endswith: \reg.exe
  - OriginalFileName: reg.exe
  selection_cli:
    CommandLine|contains|all:
    - add
    - \Software\Microsoft\Windows\CurrentVersion\Run
  condition: all of selection_*
falsepositives:
- Unknown
level: medium
title: '[AUTO] Detects the use of reg.exe to modify the "Run" registry key.'
